# build/crontab
*/5 * * * * php protected/t4.php mail/send
